# inkslab.service
[Unit]
Description=InkSlab - e-Ink TCG Card Display
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
User=root
ExecStartPre=/bin/bash /home/pi/4inch_e-Paper_E/RaspberryPi_JetsonNano/python/examples/inkslab-eink-tcg-display/scripts/selfheal.sh
ExecStart=/usr/bin/python3 -u /home/pi/4inch_e-Paper_E/RaspberryPi_JetsonNano/python/examples/inkslab-eink-tcg-display/inkslab.py
WorkingDirectory=/home/pi/4inch_e-Paper_E/RaspberryPi_JetsonNano/python/examples/inkslab-eink-tcg-display
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
