# inkslab_web.service
[Unit]
Description=InkSlab Web Dashboard
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/gunicorn --workers 1 --threads 4 --bind 0.0.0.0:80 --timeout 120 --preload --access-logfile - --error-logfile - inkslab_web:app
WorkingDirectory=/home/pi/inkslab
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
