# inkslab_web.service
[Unit]
Description=InkSlab Web Dashboard
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/waitress-serve --host=0.0.0.0 --port=80 --threads=4 --channel-timeout=120 inkslab_web:app
WorkingDirectory=/home/pi/inkslab
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
