# inkslab.service
[Unit]
Description=InkSlab - e-Ink TCG Card Display
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
User=root
ExecStartPre=/bin/bash /home/pi/inkslab/scripts/selfheal.sh
ExecStart=/usr/bin/python3 -u /home/pi/inkslab/inkslab.py
WorkingDirectory=/home/pi/inkslab
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
